# fedora/iot/startup.service
[Unit]
Description=Spark service

[Service]
ExecStart=/etc/startup-script.sh

[Install]
WantedBy=default.target